# templates/service.ini.erb
# This file is autogenerated by Puppet. Manual changes will be overwritten!
[program:<%= name -%>]
command=<%= command %>
<% if numprocs.to_i > 1 -%>
process_name=%(program_name)s_%(process_num)02d
<% end -%>
numprocs=<%= numprocs %>
priority=<%= priority %>
autostart=<%= autostart %>
autorestart=<%= autorestart %>
startsecs=<%= startsecs %>
startretries=<%= retries %>
exitcodes=<%= exitcodes %>
stopsignal=<%= stopsignal %>
stopwaitsecs=<%= stopwait %>
<% unless user == "" -%>
user=<%= user %>
<% end -%>
redirect_stderr=<%= redirect_stderr %>
<% if stdout_logfile == "" -%>
stdout_logfile=/var/log/supervisor/<%= name -%>/<%= name -%>.out
<% else -%>
stdout_logfile=<%= stdout_logfile %>
<% end -%>
stdout_logfile_maxbytes=<%= stdout_logfile_maxsize %>
stdout_logfile_backups=<%= stdout_logfile_keep %>
<% if stderr_logfile == "" -%>
stderr_logfile=/var/log/supervisor/<%= name -%>/<%= name -%>.err
<% else -%>
stderr_logfile=<%= stderr_logfile %>
<% end -%>
stderr_logfile_maxbytes=<%= stderr_logfile_maxsize %>
stderr_logfile_backups=<%= stderr_logfile_keep %>
<% unless environment == "" -%>
environment=<%= environment %>
<% end -%>
<% unless chdir == "" -%>
directory=<%= chdir %>
<% end -%>
<% unless umask == "" %>
umask=<%= umask %>
<% end -%>
